# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; ── nginx: reverse proxy on :8080 ─────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Next.js: standalone Node.js server on :3000 (internal only) ───────────
[program:nextjs]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=20
; HOSTNAME=127.0.0.1 binds Next.js to loopback — only nginx can reach it
environment=PORT="3000",HOSTNAME="127.0.0.1",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── FastAPI: uvicorn on :8000 (internal only) ─────────────────────────────
[program:fastapi]
command=uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app/backend
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
